# chef/cookbooks/baseconfig/files/default/server.conf
[program:server1]
command=/usr/bin/node /home/ubuntu/project/TeamUp/bin/www -p 8080


[group:servers]
programs=server1